# systemd/hd.mount
[Unit]
Description=Mount HD dir

[Mount]
What=/dev/sdb2
Where=/hd
Type=ntfs
Options=defaults

[Install]
WantedBy=multi-user.target
